# .config/systemd/user/battery-threshold-notify.service
[Unit]
Description=Battery Notification Script

[Service]
ExecStart=/home/donkey/bin/systemd/battery-threshold-notify.sh
